# crontab
crontab -e
* * * * * /var/backup-script
echo "test" | sendmail zkarnauhov@gmail.com
0 8 2 * * /var/backup-script
